# docs/launchd/launchd_cleanup_l1.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd 排程：每日 03:30 清理 L1 快取過期記錄

  安裝步驟：
    1. cp docs/launchd_cleanup_l1.plist.example ~/Library/LaunchAgents/com.hermes.cleanup_l1.plist
    2. launchctl load ~/Library/LaunchAgents/com.hermes.cleanup_l1.plist
    3. 驗證：launchctl list | grep hermes

  手動觸發：
    launchctl start com.hermes.cleanup_l1
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.cleanup_l1</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zhanqiru/.venvs/hermes-bio-memory/bin/python</string>
        <string>/Volumes/NO NAME/bio_DB/scheduler/cleanup_l1_cache.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zhanqiru/bio_db_backups/cleanup_l1.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zhanqiru/bio_db_backups/cleanup_l1_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
