# docs/launchd_rebuild_hnsw.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd 排程：每週日 03:00 重建 HNSW 索引

  DuckDB VSS 不支援 incremental index update，需要定期 DROP + CREATE
  才能讓新資料進入語意搜尋索引。

  安裝步驟：
    1. cp docs/launchd_rebuild_hnsw.plist.example ~/Library/LaunchAgents/com.hermes.rebuild_hnsw.plist
    2. launchctl load ~/Library/LaunchAgents/com.hermes.rebuild_hnsw.plist

  手動觸發：
    launchctl start com.hermes.rebuild_hnsw
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.rebuild_hnsw</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zhanqiru/.venvs/hermes-bio-memory/bin/python</string>
        <string>/Volumes/NO NAME/bio_DB/scheduler/rebuild_hnsw.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zhanqiru/bio_db_backups/rebuild_hnsw.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zhanqiru/bio_db_backups/rebuild_hnsw_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
