# docs/launchd_cleanup_figure_cache.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd 排程：每日 03:35 清理 figure cache 過期圖檔（TTL 預設 14 天）

  與 cleanup_l1 (03:30) 錯開 5 分鐘，避免同時 IO。
  圖檔為 content-addressed 副本，過期刪除後報告重跑會自動重建。

  安裝步驟：
    1. cp docs/launchd_cleanup_figure_cache.plist.example ~/Library/LaunchAgents/com.hermes.cleanup_figcache.plist
    2. launchctl load ~/Library/LaunchAgents/com.hermes.cleanup_figcache.plist
    3. 驗證：launchctl list | grep hermes

  手動觸發：
    launchctl start com.hermes.cleanup_figcache
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.cleanup_figcache</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zhanqiru/.venvs/hermes-bio-memory/bin/python</string>
        <string>/Volumes/NO NAME/bio_DB/scheduler/cleanup_figure_cache.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>35</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zhanqiru/bio_db_backups/cleanup_figcache.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zhanqiru/bio_db_backups/cleanup_figcache_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
